# luci-app-service-bookmarks/root/etc/uci-defaults/luci-app-service-bookmarks
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@bookmarks[-1]
	add ucitrack bookmarks
	set ucitrack.@bookmarks[-1].init=bookmarks
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
